# .config/autostart/plank.desktop
[Desktop Entry]
Type=Application
Exec=plank
Hidden=false
NoDisplay=false
Name=plank
Comment=plank
Terminal=False
